# system/etc/systemd/system/systemd-networkd.service.d/override.conf
[Unit]
After=wpa_supplicant@wlan0.service
BindsTo=wpa_supplicant@wlan0.service
